# client/rc.local
#path: /etc/rc.local

/sbin/ip link set wlan0 down

/sbin/iptables -P FORWARD ACCEPT
/sbin/iptables --table nat -A POSTROUTING -o br0 -j MASQUERADE
echo 1| tee /proc/sys/net/ipv4/ip_forward

/usr/sbin/hostapd -B /etc/hostapd/hostapd.conf &
/usr/bin/autossh -M 0 -o Tunnel=ethernet -w 0:0 -NTCf -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes root@server_ip
